# cookbooks/travis_java/templates/ubuntu/etc/profile.d/java_home.sh.erb
<% unless node['travis_java']['default_version'].empty? -%>
export JAVA_HOME=<%= node['travis_java']['jvm_base_dir'] -%>/<%= node['travis_java'][node['travis_java']['default_version']]['jvm_name'] -%>
<% end -%>
