# ci_environment/java/files/default/etc/profile.d/malloc_tune.sh

# tune glibc memory allocation
export MALLOC_ARENA_MAX=2
