# vagrant_base/redis/templates/default/redis.init.erb
#!/bin/sh
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

### BEGIN INIT INFO
# Provides:          redis
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts redis database system
# Description:       starts redis using basic start scripts
### END INIT INFO

PATH=/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC=redis-<%= node[:redis][:version] %>
DAEMON=/usr/local/bin/redis-server

test -x $DAEMON || exit 0

set -e

start_server() {
  echo -n "Starting $DESC... "
  start-stop-daemon --start --user redis -c redis:redis --startas $DAEMON -- <%= node[:redis][:config] %>
  echo "done."
}

stop_server() {
  echo -n "Stopping $DESC... "
  start-stop-daemon --stop --exec $DAEMON -c redis:redis $DAEMON -- <%= node[:redis][:config] %>
  echo "done."
}

case "$1" in
start)
  start_server
  ;;
stop)
  stop_server
  ;;
restart)
  stop_server
  sleep 4
  start_server
  ;;
*)
  N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart}" >&2
  exit 1
  ;;
esac

exit 0
