# cookbooks/travis_build_environment/templates/default/etc/init.d/disable-apparmor.sh.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          disable-apparmor
# X-Start-Before:    ssh sshd
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Disable all AppArmor profiles at boot
# Description:       Disable all AppArmor profiles at boot
### END INIT INFO

# Managed by Chef on <%= node.name %> :heart_eyes_cat:
service apparmor stop
update-rc.d -f apparmor remove
service apparmor teardown
