# ci_environment/java/templates/ubuntu/etc/profile.d/java_home.sh.erb
<% if node.java.java_home -%>
export JAVA_HOME=<%= node.java.java_home %>
<% end -%>
