# community-cookbooks/mongodb/templates/default/debian-mongodb.upstart.erb
#
# Automatically Generated by Chef, do not edit directly!
#

<%- @ulimit.to_hash.keys.sort.each do |key| %>
    <%- if not @ulimit[key].nil? %>
limit <%= key %> <%= @ulimit[key] %> <%= @ulimit[key] %>
    <% end %>
<%- end %>

kill timeout 300 # wait 300s between SIGTERM and SIGKILL.

start on runlevel [2345]
stop on runlevel [06]

script
  NAME=<%= @provides %>
  ENABLE_MONGOD="yes"
  if [ -f <%= @sysconfig_file %> ]; then
    . <%= @sysconfig_file %>;
  fi
  if [ "x$ENABLE_MONGOD" = "xyes" ]; then
  exec start-stop-daemon --start --quiet --chuid $DAEMON_USER --exec $DAEMON -- $DAEMON_OPTS
  fi
end script
