# cookbooks/travis_phantomjs/files/ubuntu/etc/profile.d/phantomjs.sh
PATH="/usr/local/phantomjs/bin:/usr/local/phantomjs:$PATH"
